# mysqldump.cron
* * * * * sudo su root /var/www/vagrant-shell-scripts/mysqldump.sh
